# deploy/gmvmax-monitor.timer
[Unit]
# Auto-deployed via .github/workflows/deploy-systemd.yml since 2026-04-11.
Description=GMVMAX Monitor Timer (every 30min, offset +15min)

[Timer]
# Fire 15 min after hour/half-hour boundaries. Two constraints:
#   1. Eventually-consistent API data: the :00 and :30 marks are common
#      GMVMAX / Shop counter refresh boundaries, so a trigger there
#      sees stale numbers and produces phantom ROI swings. Empirically
#      validated against NAD+ surge 2026-04-10 (108v/0 orders at the
#      hour edge backfilled to 397v/6 orders in the next window).
#   2. No overlap with hourly-jobs.timer at *:58 — ban_alert + ad_report
#      + balance_alert take ~90-150s and peak ~150MB on a 225MB-available
#      VPS. Keeping gmvmax-monitor at :15/:45: the :45 run finishes by
#      ~:46 (12-min gap before :58), and :15 fires 15+ min after
#      hourly-jobs winds down (~:59:30). Plenty of memory breathing room.
# :15/:45 is strictly better than :05/:35 for the API-boundary reason
# and was picked over :10/:40 because it's equidistant from :00/:30.
OnCalendar=*:15,45:00
Persistent=true

[Install]
WantedBy=timers.target
