# deploy/hourly-jobs.service
[Unit]
Description=lark-bot hourly jobs (ban_alert + ad_report + balance_alert)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=shining
WorkingDirectory=/home/shining/services/lark-bot
EnvironmentFile=/home/shining/services/lark-bot/.env
# Run ban_alert → ad_report → balance_alert in a single Python process.
# Previously this was three subprocesses behind `bash -c 'set +e; …'`,
# each re-running bootstrap.init() (~30s bitable load + tz cache + mcp
# client). Single-process reuses bootstrap state and the in-memory
# ban_status_cache singleton (so ban_alert → ad_report sees fresh bans
# without disk roundtrip), saving ~60s per timer firing.
#
# Order matches the previous bash chain. Each step is still error-
# isolated by run_with_error_fallback inside run.py — a failure pushes
# an alert card to Shining Inbox and the next step continues. run.py
# exits non-zero only if at least one step failed (better signal than
# the old `exit 0` which hid all errors from systemd).
ExecStart=/home/shining/services/lark-bot/.venv/bin/python /home/shining/services/lark-bot/run.py hourly
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hourly-jobs
TimeoutStartSec=600
# OOM guard: kill this job rather than let it starve sshd + lark-bot-ws.
# Single process holds three jobs' state; bumped 250M → 320M to give
# headroom for ad_report (~180M peak alone) plus ban_alert's discovery
# cache + balance_alert's snapshot map living in the same heap.
MemoryMax=320M
NoNewPrivileges=true
