# deploy/lark-bot-script-collect.timer
[Unit]
Description=Daily script collect timer (09:30 UTC = 17:30 CST, 30min before material-report)

[Timer]
OnCalendar=*-*-* 09:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
