# deploy/material-report.timer
[Unit]
Description=Daily Material Report Timer (18:00 CST)

[Timer]
OnCalendar=*-*-* 10:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
