# deploy/refresh-tokens.timer
[Unit]
Description=Refresh TikTok Shop Tokens Daily (17:00 CST / 09:00 UTC)

[Timer]
OnCalendar=*-*-* 09:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
