# deploy/team-lead-report-preview.service
[Unit]
Description=Team Lead Daily Report V1 Preview (17:00 CST, early signal)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=shining
WorkingDirectory=/home/shining/services/lark-bot
# V1 preview fires at 17:00 CST = 09:00 UTC = 05:00 ET for the YESTERDAY
# ET day (same date as the 18:20 final run — both report yesterday's
# complete US day). Leads get a head-start to review editor_map gaps and
# 自然流 Bitable filling before the 18:20 final push closes the day.
# During refactor grayscale, route everything to Claude ai测试 group.
ExecStart=/bin/bash -c "YESTERDAY=$(TZ=America/New_York date -d yesterday +%%Y-%%m-%%d) && exec env KANBAN_PUSH_TARGET=test .venv/bin/python run.py daily_report --view lead --date $YESTERDAY --preview --push-to-leads"
EnvironmentFile=/home/shining/services/lark-bot/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=team-lead-report-preview
MemoryMax=200M
NoNewPrivileges=true
