# deploy/team-lead-report.timer
[Unit]
Description=Team Lead Daily Report Timer (18:20 CST)

[Timer]
# 18:20 CST = 10:20 UTC. User restored original slot on 2026-04-14 after
# feedback that 18:40 felt too late. Known risk: gmvmax-monitor fires at
# *:15 and this job lands 5 min later on a 225MB-available VPS. If OOM
# cascades return (2026-04-12-style), move this to 10:05 UTC or push
# gmvmax-monitor off :15.
OnCalendar=*-*-* 10:20:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
